# vm_manager.service
[Unit]
Description=VM_Manager Flask Application
After=network.target libvirtd.service
Requires=libvirtd.service

[Service]
Type=simple
User=iris
Group=iris
WorkingDirectory=/home/iris/sisr/vm_manager/backend
Environment="FLASK_ENV=production"
Environment="VAGRANT_HOME=/home/iris/.vagrant.d"
Environment="PATH=/home/iris/sisr/vm_manager/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Commande de démarrage avec virtualenv
ExecStart=/home/iris/sisr/vm_manager/.venv/bin/gunicorn \
    --workers 3 \
    --bind 127.0.0.1:5000 \
    --timeout 300 \
    --access-logfile /var/log/vm_manager/access.log \
    --error-logfile /var/log/vm_manager/error.log \
    main:app

# Redémarrage automatique en cas d'échec
Restart=on-failure
RestartSec=5s

# Sécurité
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
