# backend/vm_manager.service
[Unit]
Description=VM_Manager Flask application (VM Manager)
After=network.target

[Service]
Type=simple
# Change User/Group if you want to run under a different account
User=edib
Group=edib
# Travailler depuis le dossier backend pour que les imports relatifs fonctionnent
WorkingDirectory=/home/edib/Vm_Manager/backend

# Environment variables (adjust if needed)
Environment=FLASK_ENV=production

# Preferred: run with gunicorn inside a virtualenv if available
# Lancer gunicorn depuis le dossier backend (import config fonctionne)
ExecStart=/bin/bash -lc 'if [ -f /home/edib/Vm_Manager/.venv/bin/activate ]; then . /home/edib/Vm_Manager/.venv/bin/activate; fi; cd /home/edib/Vm_Manager/backend && exec gunicorn --workers 3 --bind 0.0.0.0:5000 main:app'

# Fallback (uncomment if you prefer to run python directly):
# ExecStart=/usr/bin/python3 /home/edib/Vm_Manager/backend/main.py

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
